# dockerfiles/supervisord.ini
[supervisord]
nodaemon = true

[program:laravel-worker]
process_name = %(program_name)s_%(process_num)02d
command = php artisan queue:work
autostart = true
autorestart = true
user = root
numprocs = 1
redirect_stderr = true
stdout_logfile = /var/www/html/storage/logs/supervisord.log

[program:laravel-scheduler]
process_name = %(program_name)s_%(process_num)02d
command = php artisan schedule:run
autostart = true
autorestart = true
user = root
numprocs = 1
redirect_stderr = true
stdout_logfile = /var/www/html/storage/logs/supervisord.log

[program:php-fpm]
process_name = %(program_name)s_%(process_num)02d
command=php-fpm -y /usr/local/etc/php-fpm.conf -R
autostart=true
autorestart=true
user = root

[Process completed]
autostart=true
autorestart=true
numprocs=1
user=nginx
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/schedule.log
stdout_logfile_maxbytes=10MB
logfile_backups = 10
